# bin/com.johngavin.codexbar-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Install:  launchctl load -w ~/Library/LaunchAgents/com.johngavin.codexbar-refresh.plist
     Tracks:  llm#184
     Mirrors: com.johngavin.ccusage-refresh.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.johngavin.codexbar-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johngavin/docs_gh/llm/bin/refresh_codexbar_and_commit.sh</string>
    </array>

    <!-- Run every 12 hours -->
    <key>StartInterval</key>
    <integer>43200</integer>

    <!-- Do NOT run at load — first real run after manual install/test -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Environment variables needed for the script -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/johngavin/.nix-profile/bin</string>
        <key>HOME</key>
        <string>/Users/johngavin</string>
    </dict>

    <!-- Log stdout/stderr to ~/.claude/logs/ (mirrors ccusage convention) -->
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/codexbar_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/codexbar_launchd_error.log</string>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
