# .claude/launchd/com.claude.launchd-health-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.launchd-health-weekly</string>

    <!-- Run bin/launchd_health_weekly_cron.sh every Sunday at 09:00 local time -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johngavin/docs_gh/llm/bin/launchd_health_weekly_cron.sh</string>
    </array>

    <!-- Sunday = Weekday 0; Hour 9; Minute 0 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/launchd_health_weekly.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/launchd_health_weekly.err</string>

    <key>TimeOut</key>
    <integer>300</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>REPO_ROOT</key>
        <string>/Users/johngavin/docs_gh/llm</string>
    </dict>
</dict>
</plist>
