# .claude/launchd/com.claude.wiki-health-pulse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.wiki-health-pulse</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/wiki_health_check.sh</string>
        <string>/Users/johngavin/docs_gh/llm/knowledge/wiki</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/wiki_health_pulse.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/wiki_health_pulse.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
