# .claude/launchd/com.claude.config-digest-email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- com.claude.config-digest-email.plist
     Daily config-change digest: git log → markdown → email.
     Runs at 08:00 local time (same slot as the roborev daily email).

     Install:
       cp .claude/launchd/com.claude.config-digest-email.plist \
          ~/Library/LaunchAgents/com.claude.config-digest-email.plist
       launchctl load -w ~/Library/LaunchAgents/com.claude.config-digest-email.plist

     Uninstall:
       launchctl unload ~/Library/LaunchAgents/com.claude.config-digest-email.plist
       rm ~/Library/LaunchAgents/com.claude.config-digest-email.plist

     Manual dry-run:
       EMAIL_DRY_RUN=1 bash ~/docs_gh/llm/bin/config_digest_cron.sh

     Credentials: ~/.claude/env/roborev_email.env (shared with roborev email job)
       GMAIL_USERNAME=you@gmail.com
       GMAIL_APP_PASSWORD=<16-char app password>
       REPORT_RECIPIENT=you@gmail.com

     Tracked in llm#297.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.config-digest-email</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/johngavin/.local/bin/with-secrets</string>
        <string>/bin/bash</string>
        <string>/Users/johngavin/docs_gh/llm/bin/config_digest_cron.sh</string>
    </array>

    <!-- Daily at 08:05 local time (5 min after roborev daily email slot) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <!-- Do NOT run immediately on launchctl load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- PATH must include nix stable profile so nix-shell resolves -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/nix/var/nix/profiles/default/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/johngavin/.nix-profile/bin</string>
        <key>HOME</key>
        <string>/Users/johngavin</string>
        <key>NIX_SSL_CERT_FILE</key>
        <string>/nix/var/nix/profiles/default/etc/ssl/certs/ca-bundle.crt</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/config_digest_email_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/config_digest_email_launchd.log</string>

    <!-- Retry once after 60 s if job exits non-zero -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
