# .claude/launchd/com.claude.roborev-severity-autoclose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.roborev-severity-autoclose</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/roborev_severity_autoclose.sh</string>
        <string>--apply</string>
        <string>--threshold</string>
        <string>medium</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_severity_autoclose.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_severity_autoclose.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
