# .claude/launchd/com.claude.roborev-autoclose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.roborev-autoclose</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/roborev_weekly_chain.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>THRESHOLD_DAYS</key>
        <string>30</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <!-- StandardInPath=/dev/null: same SIGTTIN/state-T fix as roborev-poll-merges
         (llm#153); roborev_weekly_chain.sh calls roborev_autoclose.sh which
         uses < <(...) process subst. -->
    <key>StandardInPath</key>
    <string>/dev/null</string>
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_autoclose.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_autoclose.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
