# .claude/launchd/com.claude.worktree-gc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Sweep stale worktrees across all repos under ~/docs_gh daily 09:00.     -->
    <!-- After the script's SOAK_END date (2026-06-02), --apply removes          -->
    <!-- merged+clean+old worktrees. Before then, treated as dry-run.            -->
    <!-- Tracks:        JohnGavin/llm#199.                                       -->
    <!-- Schedule fix:  JohnGavin/llm#586 (was weekly Sun 04:00 — Mac asleep,    -->
    <!--                job never auto-fired; moved to daily 09:00 when Mac is   -->
    <!--                reliably awake).                                          -->
    <!--                                                                          -->
    <!-- Install / reload:                                                       -->
    <!--   launchctl bootout gui/$(id -u)/com.claude.worktree-gc 2>/dev/null     -->
    <!--   cp .claude/launchd/com.claude.worktree-gc.plist                       -->
    <!--        ~/Library/LaunchAgents/                                           -->
    <!--   launchctl bootstrap gui/$(id -u)                                      -->
    <!--        ~/Library/LaunchAgents/com.claude.worktree-gc.plist              -->

    <key>Label</key>
    <string>com.claude.worktree-gc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/worktree_gc.sh</string>
        <string>--apply</string>
    </array>

    <!-- Daily 09:00 (Mac reliably awake) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/worktree_gc.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/worktree_gc.err</string>

    <!-- Tolerate up to 30 s between retries if a minute is missed -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Do NOT run on load — only on the schedule. -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
