# .claude/launchd/com.claude.roborev-severity-autoclose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.roborev-severity-autoclose</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/roborev_severity_autoclose.sh</string>
        <string>--apply</string>
        <string>--threshold</string>
        <string>medium</string>
    </array>
    <!-- ~/.local/bin/ MUST precede /usr/local/bin/ so the Phase-1.6 roborev
         primary shim (#386) intercepts all roborev invocations from severity-autoclose. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/johngavin/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <!-- StandardInPath=/dev/null: same SIGTTIN/state-T fix as roborev-poll-merges
         (llm#153); roborev_severity_autoclose.sh uses < <(...) process subst. -->
    <key>StandardInPath</key>
    <string>/dev/null</string>
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_severity_autoclose.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_severity_autoclose.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
