# .claude/launchd/com.claude.roborev-poll-merges.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.roborev-poll-merges</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/johngavin/docs_gh/llm/.claude/scripts/roborev_poll_merges.sh</string>
        <string>--apply</string>
    </array>
    <!-- Business-hours safety-net poller (Phase 1.7, #217):
         Fires thrice per weekday at 09:00, 13:00, 17:00 — 3 runs/day × 5 days = 15/week.
         Primary coverage now comes from the post-merge git hook installed per-repo via
         roborev_install_post_merge_hook.sh; the poller is the fallback for repos that
         don't have the hook yet, or for merges that bypass the local git workflow.
         Prior schedule was hourly 09-21 (65/week). Before that, every 15 min 24/7.
         Phase 4 (full poller removal) deferred pending 7-day soak after hook rollout. -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Tuesday -->
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Wednesday -->
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Thursday -->
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Friday -->
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <!-- ~/.local/bin/ MUST precede /usr/local/bin/ so the Phase-1.6 roborev
         primary shim (#386) intercepts all roborev invocations from poll-merges. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/johngavin/.local/bin:/Users/johngavin/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <!--
      StandardInPath=/dev/null prevents SIGTTIN/SIGTTOU stopping the process
      (state T) when bash's process substitution < <(...) opens a pipe under
      bash 3.2 without a controlling terminal. Root cause: launchd by default
      provides an open stdin fd; bash job-control interprets a read attempt on
      that fd as a background tty-read and sends SIGTSTP/SIGTTIN → state T.
      Setting stdin to /dev/null severs that pathway. Fixes llm#153.
    -->
    <key>StandardInPath</key>
    <string>/dev/null</string>
    <key>StandardOutPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_poll_merges.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/johngavin/.claude/logs/roborev_poll_merges.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
